# files/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:rails]
command=/bin/bash -c "cd /myapp && /usr/local/bin/bundle exec rails server"
