# files/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:rails]
command=/bin/bash -c "cd /myapp && ./bin/rails server"

[program:delayed_job]
command=/bin/bash -c "cd /myapp && ./bin/delayed_job run"
